# services/spec-preview-server/systemd/spec-preview-server-update.timer
[Unit]
Description=Periodically run spec-preview-server self-update

[Timer]
OnBootSec=2min
OnUnitActiveSec=10min
Persistent=true

[Install]
WantedBy=timers.target
